# petsona.service
[Unit]
Description=Petsona Flask Application
After=network.target
Requires=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/petsona/app
Environment="PATH=/path/to/your/venv/bin"
Environment="FLASK_ENV=production"
Environment="SECRET_KEY=your-super-secret-key-here"
Environment="DATABASE_URI=mysql+pymysql://user:pass@host/db"
ExecStart=/path/to/your/venv/bin/gunicorn \
    --worker-class eventlet \
    --workers 1 \
    --worker-connections 1000 \
    --bind unix:/tmp/petsona.sock \
    --access-logfile /var/log/petsona/access.log \
    --error-logfile /var/log/petsona/error.log \
    --log-level info \
    --capture-output \
    --enable-stdio-inheritance \
    run:app

# Restart on failure
Restart=always
RestartSec=5

# Security
NoNewPrivileges=yes
PrivateTmp=yes

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target